# debian/wintersky-netinstall-livesettings.service
[Unit]
Description=netinstall cheatcodes in live media
ConditionPathExists=/usr/share/wintersky-netinstall-livesettings/
After=wintersky-common-livesettings.service

[Service]
Type=oneshot
ExecStart=/usr/share/wintersky-netinstall-livesettings/netinstall-cheatcodes start

[Install]
WantedBy=graphical.target
